# opt/stunnel/.footprint
drwxr-xr-x etc/
drwxr-xr-x etc/rc.d/
-rwxr-xr-x etc/rc.d/stunnel
drwxr-xr-x etc/stunnel/
-rw-r--r-- etc/stunnel/stunnel.conf-sample
drwxr-xr-x usr/
drwxr-xr-x usr/bin/
-rwxr-xr-x usr/bin/stunnel
-rwxr-xr-x usr/bin/stunnel3
drwxr-xr-x usr/lib/
drwxr-xr-x usr/lib/stunnel/
-rwxr-xr-x usr/lib/stunnel/libstunnel.so
drwxr-xr-x usr/share/
drwxr-xr-x usr/share/man/
drwxr-xr-x usr/share/man/man8/
-rw-r--r-- usr/share/man/man8/stunnel.8
-rw-r--r-- usr/share/man/man8/stunnel.fr.8
-rw-r--r-- usr/share/man/man8/stunnel.pl.8
drwxr-xr-x usr/var/
drwxr-xr-x usr/var/lib/
drwxrwx--T usr/var/lib/stunnel/
